# s6/services.d/annatar-workers/run
#!/bin/bash

with-contenv rq worker-pool -n ${WORKERS:-5} \
	annatar:${NAMESPACE:-annatar}:search \
	annatar:${NAMESPACE:-annatar}:curator
